# Features/server-optimization/service.service
[Unit]
Description=%N

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /usr/local/bin/%N.sh

[Install]
WantedBy=multi-user.target 